# systemd/.config/systemd/user/backup.service
[Unit]
Description=Cloud backups daemon

[Service]
ExecStart=/home/simka/.bin/backup
Restart=always

[Install]
WantedBy=default.target
